# deploy/systemd/nakshatra-unconscious.service.d/70-from-roster.conf
# from_roster cutover (slice 4 capstone): the unconscious chain is GENERATED from the live roster
# (~/.nakshatra/unconscious-roster.tsv), firewall-gated, slices self-provisioned from the package,
# autonomously summoned + scale-to-zero. ROLLBACK: rm this file; systemctl --user daemon-reload;
# systemctl --user restart nakshatra-unconscious.service  → reverts to the static chain_yaml + @a/@b units.
[Service]
Environment=MODELS=%h/.nakshatra/serve_models.unconscious-roster.yaml
Environment=MESH_PEERS=%h/.nakshatra/unconscious-roster.tsv
# disable the static @a/@b unit controller — the roster controller summons now
Environment="NAKSHATRA_LIFECYCLE_UNITS="
Environment="NAKSHATRA_LIFECYCLE_PROBES="
Environment=NAKSHATRA_LIFECYCLE_ROSTER_MODEL=prithvi-unconscious
Environment=NAKSHATRA_LIFECYCLE_ROSTER_HIDDEN_SIZE=4096
Environment=NAKSHATRA_LIFECYCLE_ROSTER_NUM_LAYERS=32
Environment=NAKSHATRA_LIFECYCLE_ROSTER_PACKAGE=%h/.nakshatra/packages/dsr1-llama8b
Environment=NAKSHATRA_LIFECYCLE_DAEMON_BIN=%h/llama.cpp/build/bin/llama-nakshatra-worker
Environment=NAKSHATRA_LIFECYCLE_PYTHON_BIN=%h/nakshatra/.venv/bin/python
# 2026-06-18: clear the stale Sthambha pillar lease inherited from the base unit (dead model-id +
# down pillar → the reaper's fail-safe pinned the chain forever). Single box → local idle-grace governs.
Environment=NAKSHATRA_LIFECYCLE_PILLAR_URL=
Environment=NAKSHATRA_LIFECYCLE_LEASE_MODEL=
Environment=NAKSHATRA_LIFECYCLE_LEASE_KEYID=
Environment=NAKSHATRA_LIFECYCLE_IDLE_GRACE_S=120
