# deploy/systemd/nakshatra-unconscious-worker@.service
[Unit]
Description=Nakshatra unconscious-tier worker %i (layer-split gRPC worker)
Documentation=https://github.com/fthrvi/nakshatra/blob/main/docs/INFRA-MAP.md
After=network.target

[Service]
Type=simple
# Per-instance config: ~/.nakshatra/unconscious-worker-<a|b>.env defines
#   PORT, SUB_GGUF, MODE, LAYER_START, LAYER_END, MODEL_ID, N_CTX, N_GPU_LAYERS
EnvironmentFile=%h/.nakshatra/unconscious-worker-%i.env
Environment=DAEMON_BIN=%h/llama.cpp/build/bin/llama-nakshatra-worker
Environment=MODEL_ID=nakshatra-unconscious
Environment=N_CTX=2048
Environment=N_GPU_LAYERS=99
# L3 (Sthambha) registration — set PILLAR_URL in the env file to register this
# worker as a compute peer (self-registers + heartbeats every 30s). Empty → no
# registration (legacy static-chain behaviour). GPU_* describe the offered HW.
Environment=PILLAR_URL=
Environment=GPU_VENDOR=AMD
Environment=GPU_MODEL=RX-9070-XT
Environment=GPU_VRAM_GB=16
Environment=GPU_BACKEND=rocm
Environment=VRAM_OFFERED_GB=6
ExecStart=%h/nakshatra/.venv/bin/python %h/nakshatra/scripts/worker.py \
    --port ${PORT} --sub-gguf ${SUB_GGUF} --mode ${MODE} \
    --layer-start ${LAYER_START} --layer-end ${LAYER_END} \
    --model-id ${MODEL_ID} --daemon-bin ${DAEMON_BIN} \
    --n-ctx ${N_CTX} --n-gpu-layers ${N_GPU_LAYERS} \
    --pillar-url=${PILLAR_URL} --node-id=unconscious-%i \
    --gpu-vendor=${GPU_VENDOR} --gpu-model=${GPU_MODEL} \
    --gpu-vram-gb=${GPU_VRAM_GB} --gpu-backend=${GPU_BACKEND} \
    --vram-offered-gb=${VRAM_OFFERED_GB}
Restart=always
RestartSec=4

[Install]
WantedBy=default.target
